# templates/systemd/dehydrated.service
[Service]
ExecStart={{ dehydrated.prefix.bin }}/dehydrated --accept-terms --cron
