# templates/systemd/dehydrated.timer
[Timer]
{{ dehydrated.systemd.timer -}}
Persistent=true

[Install]
WantedBy=timers.target
