# templates/systemd/dehydrated.service
[Service]
ExecStart=dehydrated --accept-terms --cron
