# templates/systemd/dehydrated.service
[Service]
Type=oneshot
ExecStart={{ dehydrated.prefix.bin }}/dehydrated --accept-terms --cron
